# files/traefik/whoami.container
[Unit]
Description=whoami Container
Requires=traefik.network
After=traefik.network

[Container]
ContainerName=whoami
Image=docker.io/traefik/whoami:latest

SecurityLabelDisable=true

AutoUpdate=registry

Network=traefik.network

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target
